# Module 1/Task2_Backup/backup.service
[Unit]
Description=Automated Backup Service
Wants=backup.timer

[Service]
User=root
ExecStart=/usr/local/bin/backup.sh

